# etc/init/xe-automator.conf
# xe automator

description     "XenServer Guest Configuration Automator"

start on (filesystem
          and started networking)

script
        /usr/sbin/xe-set-network
        /usr/sbin/xe-set-hostname
        /usr/sbin/generate-sshd-keys
end script
